# services/powertop.timer
[Unit]
Description=run powertop autotune every 10 minutes

[Timer]
OnBootSec=15min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target

